# control/prerm
#!/bin/sh

/etc/init.d/videochat stop > /dev/null 2>&1

update-rc.d -f videochat remove > /dev/null 2>&1

exit 0
